# com.dna-commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dna-commit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/funahashihotaka/dev/DNA-commit/venv/bin/python</string>
        <string>/Users/funahashihotaka/dev/DNA-commit/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/funahashihotaka/dev/DNA-commit</string>

    <!-- 複数の時間帯に実行（スリープ対策） -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 早朝 3:00 -->
        <dict>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 朝 9:00（Macが起動している可能性が高い） -->
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 夜 21:00（作業後） -->
        <dict>
            <key>Hour</key>
            <integer>21</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/funahashihotaka/dev/DNA-commit/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/funahashihotaka/dev/DNA-commit/logs/launchd_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/funahashihotaka/dev/DNA-commit/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- ロード時に一度実行（Macの起動/再起動後に確実に実行） -->
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
