# com.dna-commit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dna-commit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/funahashihotaka/dev/DNA-commit/venv/bin/python</string>
        <string>/Users/funahashihotaka/dev/DNA-commit/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/funahashihotaka/dev/DNA-commit</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/funahashihotaka/dev/DNA-commit/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/funahashihotaka/dev/DNA-commit/logs/launchd_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/funahashihotaka/dev/DNA-commit/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
